# deploy/fambank.service
[Unit]
Description=FamBank Backend API Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=fambank
Group=fambank
WorkingDirectory=/opt/fambank/backend
EnvironmentFile=/opt/fambank/backend/.env
ExecStart=/opt/fambank/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000

# Auto-restart on failure
Restart=on-failure
RestartSec=5

# Burst limit: max 5 restarts within 60 seconds
StartLimitBurst=5
StartLimitIntervalSec=60

# Logging to journald (structured JSON from app)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fambank

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/fambank

[Install]
WantedBy=multi-user.target
